# assets/etc/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:tractor-engine]
command=/opt/pixar/Tractor-2.2/bin/tractor-engine --configdir=/opt/pixar/config --log=/var/log/engine.log --debug --verbose
autorestart=true

[program:tractor-blade]
command=/opt/pixar/Tractor-2.2/bin/tractor-blade --log=/var/log/blade.log --debug --verbose
autorestart=true
